# etc/serverlens.service
[Unit]
Description=ServerLens MCP Server
After=network.target postgresql.service
Documentation=https://gitlab.rucode.org/devtools/serverlens

[Service]
Type=simple
User=serverlens
Group=serverlens
ExecStart=/opt/serverlens/venv/bin/python -m serverlens serve --config /etc/serverlens/config.yaml
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Environment
EnvironmentFile=-/etc/serverlens/env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/
ReadWritePaths=/var/log/serverlens
PrivateTmp=yes
CapabilityBoundingSet=
SystemCallFilter=@system-service
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
